# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:fastapi_app]
command=bash -c 'uvicorn api.main:app --host 0.0.0.0 --port ${CONTAINERPORT_API:-8000}'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=bash -c 'npm start -- -p ${CONTAINERPORT_FRONTEND:-3000}'
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0 